# packaging/archlinux/kaillux.service
[Unit]
Description=Kaillux Kaillera Server
Documentation=https://github.com/kidoz/kaillux
After=network.target

[Service]
Type=simple
User=kaillux
Group=kaillux

# Java options
Environment="KAILLUX_OPTS=-Xms64m -Xmx256m"

# Working directory for logs
WorkingDirectory=/var/lib/kaillux

# Start command
ExecStart=/usr/bin/kaillux

# Restart on failure
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Allow write to working directory and logs
ReadWritePaths=/var/lib/kaillux

# Network access
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
